# ansible/roles/minio/templates/minio-exporter.service.j2
[Unit]
Description=MinIO Exporter
Wants=network-online.target
After=network-online.target

[Service]
User=minio
Group=minio
ExecStart=/usr/local/bin/minio-exporter \
  --minio.server={{ minio_url }} \
  --minio.access-key={{ minio_service_user }} \
  --minio.secret-key={{ minio_service_password }}
Restart=always

[Install]
WantedBy=multi-user.target
